# roles/cx_container_app/templates/container.socket.j2
[Socket]
{% for stream in cx_container_app_socket_listen_streams %}
ListenStream={{ stream }}
{% endfor %}

[Install]
WantedBy=sockets.target
